# systemfiles/etc/systemd/system/post-boot-script.service

[Unit]
Description=post-boot-script

[Service]
ExecStart=/etc/systemd/system/post-boot-script.sh

[Install]
WantedBy=multi-user.target

